# backup-py.service
[Unit]
Description=backup-py

[Service]
Type=simple
ExecStart=/usr/local/bin/backup.py

[Install]
WantedBy=multi-user.target
